# raspberry-pi/piespeakup2/config/piespeakup2restart.timer
[Unit]
Description=Timed restart of piespeakup2.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
