# raspberry-pi/piespeakup2/config/piespeakup2restart.service
[Unit]
Description=Restart piespeakup2.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart piespeakup2.service
