# raspberry-pi/piespeakup2/config/piespeakup2.service.in
[Unit]
Description=Software speech output for Speakup via the GPU on the Raspberry Pi

[Service]
Type=notify
WatchdogSec=43200
PIDFile=/run/piespeakup2.pid
ExecStart=prefix/bin/piespeakup2 --device=local
ExecReload=/bin/kill -HUP $MAINPID
Restart=always


[Install]
WantedBy=multi-user.target
