# generic/speechd-up-ng/config/speechd-upd.service.in
[Unit]
Description=speechd-up, a server to connect SpeakUp to speech-dispatcher
After=syslog.target

[Service]
Type=notify
WatchdogSec=3600
PIDFile=/run/speechd-up.pid
ExecStart=prefix/bin/speechd-up -d
ExecReload=/bin/kill -HUP $MAINPID
Restart=always


[Install]
WantedBy=multi-user.target
